# aux/upstart/mechanicd.conf
description "Android-isms and kernel features on GNU/Linux systems"

start on started dbus
stop on stopping dbus

respawn

exec /usr/bin/mechanicd